# packages/spartaqube/spartaqube-0.1.73.tar.gz/spartaqube-0.1.73/spartaqube/supervisor/supervisor.conf
[supervisord]
# environment = HOST_DB_EE=%(ENV_HOST_DB_EE)s, MYSQL_DATABASE=%(ENV_MYSQL_DATABASE)s, MYSQL_USER=%(ENV_MYSQL_USER)s, MYSQL_PASSWORD=%(ENV_MYSQL_PASSWORD)s, SERVICE_DB_EE_SERVICE_HOST=%(ENV_SERVICE_DB_EE_SERVICE_HOST)s
# stdout_logfile = /app/supervisor/log/supervisord/supervisord-stdout.log
# stderr_logfile = /app/supervisor/log/supervisord/supervisord-stderr.log
stderr_logfile=/dev/null  ; Redirect standard error to /dev/null
stdout_logfile=/dev/null  ; Redirect standard output to /dev/null


[program:spartaqube-uvicorn]
# command = /usr/local/bin/gunicorn spartaqube_app.wsgi:application --bind 0.0.0.0:8000
# uvicorn_server = f"uvicorn 'spartaqube_app'.asgi:application --workers {workers} --host 0.0.0.0 --port {port} --loop uvloop --http httptools --log-level warning &"
# command = /usr/local/bin/gunicorn spartaqube_app.wsgi:application -c /app/gunicorn/gunicorn.conf.py
# command = /usr/local/bin/uvicorn spartaqube_app.asgi:application --host 0.0.0.0 --port 8665 --loop uvloop --http httptools --log-level warning
command=/bin/bash /app/supervisor/start_uvicorn.sh
user = root
directory = /usr/local/lib/python3.11/site-packages/spartaqube
autostart = true
autorestart = true
