# packages/spartaqube/spartaqube-0.1.73.tar.gz/spartaqube-0.1.73/spartaqube/supervisor/supervisorDaphneDEPREC.conf
[supervisord]
# environment = HOST_DB_EE=%(ENV_HOST_DB_EE)s, MYSQL_DATABASE=%(ENV_MYSQL_DATABASE)s, MYSQL_USER=%(ENV_MYSQL_USER)s, MYSQL_PASSWORD=%(ENV_MYSQL_PASSWORD)s, SERVICE_DB_EE_SERVICE_HOST=%(ENV_SERVICE_DB_EE_SERVICE_HOST)s
# stdout_logfile = /app/supervisor/log/supervisord/supervisord-stdout.log
# stderr_logfile = /app/supervisor/log/supervisord/supervisord-stderr.log
stderr_logfile=/dev/null  ; Redirect standard error to /dev/null
stdout_logfile=/dev/null  ; Redirect standard output to /dev/null

[program:spartaqube-daphne]
command = /usr/local/bin/daphne spartaqube_app.asgi:application -b 0.0.0.0 -p 5000
user = root
directory = /usr/local/lib/python3.11/site-packages/spartaqube
autostart = true
autorestart = true
# stdout_logfile = /app/supervisor/log/daphne/daphne-stdout.log
# stderr_logfile = /app/supervisor/log/daphne/daphne-stderr.log
stderr_logfile=/dev/null  ; Redirect standard error to /dev/null
stdout_logfile=/dev/null  ; Redirect standard output to /dev/null

[program:spartaqube-gunicorn]
# command = /usr/local/bin/gunicorn spartaqube_app.wsgi:application --bind 0.0.0.0:8000
command = /usr/local/bin/gunicorn spartaqube_app.wsgi:application -c /app/gunicorn/gunicorn.conf.py
user = root
directory = /usr/local/lib/python3.11/site-packages/spartaqube
autostart = true
autorestart = true
